# supervisor/by-mark.conf
[program:technoinst]
command=/home/mark/by-mark.ru/bin/start_gunicorn.sh
user=mark
process_name=%(program_name)s
priority=1
numproc=1
autostart=1
autorestart=1
redirect_stderr=true
